# misc/openwrt/files/etc/crontabs/root
# 0 0 * * * sleep 70 && uci set sqm.@queue[0].enabled='1' && uci commit sqm && /etc/init.d/sqm reload
0 0 * * * sleep 70 && uci set sqm.@queue[0].enabled='0' && uci commit sqm && /etc/init.d/sqm reload
